# network-tuning.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='SUNW:network_tuning'>

<service
       name='site/network_tuning'
       type='service'
       version='1'>

       <create_default_instance enabled='true' />

       <single_instance />

       <dependency
               name='usr'
               type='service'
               grouping='require_all'
               restart_on='none'>
               <service_fmri value='svc:/system/filesystem/minimal' />
       </dependency>

<!-- Run ndd commands after network/physical is plumbed. -->
       <dependency
               name='network-physical'
               grouping='require_all'
               restart_on='none'
               type='service'>
               <service_fmri value='svc:/network/physical' />
       </dependency>

<!-- but run the commands before network/initial -->
       <dependent
               name='ndd_network-initial'
               grouping='optional_all'
               restart_on='none'>
               <service_fmri value='svc:/network/initial' />
       </dependent>

       <exec_method
               type='method'
               name='start'
               exec='/lib/svc/method/network-tuning.sh'
               timeout_seconds='60' />

       <exec_method
               type='method'
               name='stop'
               exec=':true'
               timeout_seconds='60' />

       <property_group name='startd' type='framework'>
               <propval name='duration' type='astring' value='transient' />
       </property_group>

       <stability value='Unstable' />

       <template>
               <common_name>
                       <loctext xml:lang='C'>
                               Network Tunings
                       </loctext>
               </common_name>

       </template>
</service>

</service_bundle>
